# deploy/systemd/trading-upgrade.service
[Unit]
Description=Trading System Upgrade

[Service]
Type=oneshot
User=trading
Group=trading
WorkingDirectory=/opt/trading/app
EnvironmentFile=/etc/trading/trading.env
Environment=TRADING_UPGRADE_SERVICE_CONTROL=0
ExecStartPre=+/bin/systemctl stop trading-operator.service
ExecStartPre=+/bin/systemctl stop trading-engine.service
ExecStart=/opt/trading/app/deploy/bin/upgrade_trading_system.sh
ExecStartPost=+/usr/bin/install -m 0644 /opt/trading/app/deploy/logrotate/trading-system /etc/logrotate.d/trading-system
ExecStartPost=+/usr/sbin/logrotate -d /etc/logrotate.d/trading-system
ExecStartPost=+/bin/systemctl start trading-engine.service
ExecStartPost=+/bin/systemctl restart trading-operator.service
StandardOutput=append:/var/lib/trading/logs/upgrade.log
StandardError=append:/var/lib/trading/logs/upgrade.log
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/trading /var/backups/trading /opt/trading/app/data /opt/trading/app/logs /etc/trading
