# deploy/systemd/trading-restore-drill.timer
[Unit]
Description=Run trading restore drill weekly

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600
Persistent=true
Unit=trading-restore-drill.service

[Install]
WantedBy=timers.target
