# deploy/systemd/trading-backup.timer
[Unit]
Description=Trading System Daily Database Backup

[Timer]
OnCalendar=daily
Persistent=true
Unit=trading-backup.service

[Install]
WantedBy=timers.target