# deploy/systemd/trading-backup.service
[Unit]
Description=Trading System Legacy SQLite Backup

[Service]
Type=oneshot
User=trading
Group=trading
WorkingDirectory=/opt/trading/app
EnvironmentFile=/etc/trading/trading.env
ExecStart=/opt/trading/app/deploy/bin/backup_trading_db.sh
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/trading /var/backups/trading /opt/trading/app/data /opt/trading/app/logs /etc/trading
