# ops/server/systemd/trading-restore-drill.timer
[Unit]
Description=Run trading restore drill monthly

[Timer]
OnCalendar=monthly
Persistent=true
Unit=trading-restore-drill.service

[Install]
WantedBy=timers.target
