# ops/server/systemd/trading-backup-evidence.service
[Unit]
Description=Trading backup/WAL/restore evidence gate
Documentation=file:/opt/trading/app/ops/server/README.md
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=postgres
Group=trading
EnvironmentFile=-/etc/trading/trading.env
EnvironmentFile=-/etc/trading/provider.env
Environment=BACKUP_EVIDENCE_REQUIRE_SIGNATURE=1
Environment=PGHOST=/var/run/postgresql
Environment=PGPORT=5432
Environment=PGUSER=postgres
Environment=TS_BACKUP_EVIDENCE_RUN_BASE_BACKUP=0
Environment=TS_BACKUP_EVIDENCE_RUN_RESTORE_DRILL=0
Environment=TS_BACKUP_EVIDENCE_RUN_WAL_CATCHUP=0
Environment=TS_BACKUP_EVIDENCE_PROBE_TIMEOUT_S=5
Environment=TS_BACKUP_EVIDENCE_SYSTEMCTL_TIMEOUT_S=5
ExecStart=/bin/bash /opt/trading/ops/backup/backup_restore_evidence.sh
TimeoutStartSec=8min
Restart=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-backup-evidence
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/trading /tmp /var/tmp

[Install]
WantedBy=multi-user.target
