# ops/server/systemd/trading-restore-drill.service
[Unit]
Description=Trading restore drill
Documentation=file:/opt/trading/app/ops/server/README.md
Wants=network-online.target
After=network-online.target postgresql.service pgbouncer.service

[Service]
Type=oneshot
User=postgres
Group=trading
EnvironmentFile=-/etc/trading/trading.env
EnvironmentFile=-/etc/trading/provider.env
ExecStart=/bin/bash /opt/trading/ops/backup/restore_drill.sh
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-restore-drill
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/trading /tmp /var/tmp

[Install]
WantedBy=multi-user.target
